# crates/safu-node/deploy/safu-node.service
# Sample systemd unit for the always-on safu-node backup/share host.
#
# Setup:
#   sudo useradd --system --home /var/lib/safu-node --create-home safu
#   sudo install -m 0700 -o safu -g safu -d /var/lib/safu-node
#   # Put the passphrase in a root-only env file (keeps it out of the unit and
#   # the process list); SAFU_NODE_PASSPHRASE is read from the environment:
#   printf 'SAFU_NODE_PASSPHRASE=%s\n' 'your-passphrase' | sudo tee /etc/safu-node.env >/dev/null
#   sudo chmod 600 /etc/safu-node.env
#   sudo cp safu-node.service /etc/systemd/system/
#   sudo systemctl enable --now safu-node
#
# Updating is then: install the new binary (the install script), and
#   sudo systemctl restart safu-node
# The data dir is untouched, so devices stay linked.

[Unit]
Description=safu-node — zero-knowledge backup & public-share host
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=safu
Group=safu
EnvironmentFile=/etc/safu-node.env
Environment=SAFU_NODE_DATA_DIR=/var/lib/safu-node
ExecStart=/usr/local/bin/safu-node serve
# Graceful stop: SIGTERM lets the node flush its document snapshot before exit.
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=5

# Hardening: the node only needs its own data dir.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/safu-node
PrivateTmp=true

[Install]
WantedBy=multi-user.target
